# service/rclone/systemd/rclone-push@.service
[Unit]
Description=Rclone Push for /%I
Wants=container-build@rclone.service
After=container-build@rclone.service
ConditionDirectoryNotEmpty=/%I

[Service]
Type=oneshot
SyslogIdentifier=%N
Environment=DESTINATION=crypt:%i
ExecStart=/bin/podman run --rm --env-file %E/coreos-home-server/rclone/rclone.env --user root --volume /%I:/data:z localhost/rclone:latest sync /data ${DESTINATION}

[Install]
WantedBy=multi-user.target
