# service/rclone/systemd/rclone-push-volume@.timer
[Unit]
Description=Scheduled Rclone Push for Volume %I

[Timer]
OnCalendar=02:00
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
