# service/dovecot/systemd/grafana-dashboard@dovecot.service.d/after-prometheus.conf
[Unit]
Wants=prometheus-service-register@dovecot:9900.service
After=prometheus-service-register@dovecot:9900.service
