# service/prosody/systemd/grafana-dashboard@prosody.service.d/after-prometheus.conf
[Unit]
Wants=prometheus-service-register@prosody:5280.service
After=prometheus-service-register@prosody:5280.service
