# config/common/systemd/container-build@.path
[Unit]
Description=Container build for %I

[Path]
PathModified=%E/coreos-home-server/%i/Containerfile
Unit=container-build@%i.service

[Install]
WantedBy=multi-user.target
