# config/common/systemd/container-build@.path
[Unit]
Description=Container build for %I

[Path]
PathModified=/etc/container-service/%i/Containerfile
Unit=container-build@%i.service

[Install]
WantedBy=multi-user.target
