# service/biboumi/quadlet/biboumi.container
[Unit]
Description=Biboumi IRC Gateway for XMPP
Wants=container-build@%N.service container-volume@%N.service prosody.service prosody-component-register@%N.service
After=container-build@%N.service container-volume@%N.service prosody.service prosody-component-register@%N.service

[Container]
ContainerName=%N
EnvironmentFile=%E/coreos-home-server/%N/%N.env
Image=localhost/%N:latest
Network=internal
Volume=%N:/var/lib/%N:z

[Service]
Restart=on-failure

[Install]
WantedBy=multi-user.target
