# service/rclone/systemd/container-volume-restore@.service.d/before-rclone.conf
[Unit]
Wants=rclone-pull@var-lib-backups-coreos\x2dhome\x2dserver-%i.service
Before=rclone-pull@var-lib-backups-coreos\x2dhome\x2dserver-%i.service
