# host/base/systemd/container-volume-backup@.timer
[Unit]
Description=Scheduled Backup for Container Volume %i

[Timer]
OnCalendar=02:00
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
