# host/base/systemd/coreos-home-server-update.timer
[Unit]
Description=Scheduled Update for CoreOS Home Server

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
