# config/service/letsencrypt/letsencrypt-dns-renew@.timer
[Unit]
Description="Let's Encrypt DNS01 scheduled certificate renewal for %I"

[Timer]
OnUnitActiveSec=7d

[Install]
WantedBy=timers.target
