# service/prometheus/systemd/grafana-dashboard@prometheus-podman-exporter.service.d/after-prometheus.conf
[Unit]
Wants=prometheus-service-register@prometheus-podman-exporter:9882.service
After=prometheus-service-register@prometheus-podman-exporter:9882.service
ConditionPathExists=
ConditionPathExists=%E/coreos-home-server/prometheus/service/grafana-dashboard-podman-exporter.json

[Service]
Environment=DASHBOARD_CONF=%E/coreos-home-server/prometheus/service/grafana-dashboard-podman-exporter.json
