# service/prometheus/systemd/grafana-dashboard@prometheus.service.d/after-prometheus.conf
[Unit]
Wants=prometheus-service-register@host.containers.internal:9100.service
After=prometheus-service-register@host.containers.internal:9100.service
