# service/prometheus/systemd/grafana-dashboard@prometheus-node-exporter.service.d/after-prometheus.conf
[Unit]
Wants=prometheus-service-register@host.containers.internal:9100.service prometheus-service-register@prometheus-smartctl-exporter:9633.service
After=prometheus-service-register@host.containers.internal:9100.service prometheus-service-register@prometheus-smartctl-exporter:9633.service
ConditionPathExists=
ConditionPathExists=%E/coreos-home-server/prometheus/service/grafana-dashboard-node-exporter.json

[Service]
Environment=DASHBOARD_CONF=%E/coreos-home-server/prometheus/service/grafana-dashboard-node-exporter.json
